# texas_supervisor.conf
[program:texas_uwsgi]
command=/root/env/bin/uwsgi --ini /root/texas/texas_uwsgi.ini
stderr_logfile=/root/texas/logs/uwsgi_stderr.log
stderr_logfile_maxbytes=0
stopasgroup=true
