# texas_scheduler_supervisor.conf
[program:texas_scheduler]
command=/root/env/bin/python manage.py rqscheduler
directory=/root/texas
user=root
autostart=true
autorestart=true
